# Dockerfile/supervisord.ini
[supervisord]
nodaemon=true

[program:node]
command=/usr/bin/node /etc/guacamole.js &

[program:guacd]
command=/usr/local/sbin/guacd -b 0.0.0.0 -
